# ansible/roles/02-chroot/templates/firstboot.service.j2
[Unit]
Description=firstboot

[Service]
Type=oneshot
ExecStart=/opt/firstboot
ExecStop=/opt/firstboot
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
